# kalshi/kalshi-weather.service.example
[Unit]
Description=Kalshi multi-city weather trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/weatherbot
EnvironmentFile=/opt/weatherbot/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/weatherbot/.venv/bin/python /opt/weatherbot/kalshi_weather_trader.py run --config /opt/weatherbot/kalshi_weather_config.json
Restart=always
RestartSec=10
TimeoutStopSec=30
# SIGINT lets Python unwind and run the registered execution-manager cleanup.
KillSignal=SIGINT
UMask=0077
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
